# scripts/crontab
# Set this once to the absolute path where you cloned TeleVault, e.g. /home/youruser/TeleVault
TELEVAULT_HOME=/path/to/TeleVault

0 * * * * mkdir -p $TELEVAULT_HOME/logs && $TELEVAULT_HOME/scripts/backup_db.sh >> $TELEVAULT_HOME/logs/backup.log 2>&1